# get_quote.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/get_quote.sh

[Install]
WantedBy=default.target